# systemd/ccpipe.service.in
[Unit]
Description=ccpipe — Claude Code remote web interface
After=default.target ccpipe-virtual-mic.service
# Soft dep: try to bring up the virtual mic so /voice works. If it can't
# (e.g. Pulse not running), ccpipe still starts — text + TTS playback
# don't depend on it.
Wants=ccpipe-virtual-mic.service

[Service]
Type=exec
# @REPO_ROOT@ is substituted by scripts/install.sh at install time with
# the absolute path of the checkout being installed. The rendered unit
# lives in ~/.config/systemd/user/ccpipe.service; the source template
# stays in the repo.
WorkingDirectory=@REPO_ROOT@/backend
# Default ExecStart binds 0.0.0.0 so an off-host reverse proxy (nginx)
# can reach the backend over the LAN. When TLS is terminating at that
# proxy, add a drop-in that resets ExecStart and re-runs uvicorn with
# `--proxy-headers --forwarded-allow-ips=<proxy-host-IP>`; without
# the latter the per-IP login throttle treats every request as coming
# from the proxy. See README §"Reverse proxy" for the exact drop-in.
ExecStart=@REPO_ROOT@/backend/.venv/bin/uvicorn ccpipe.main:app --host 0.0.0.0 --port 8080

# Frontend bundle served by FastAPI when this exists
Environment=CCPIPE_FRONTEND_DIST=@REPO_ROOT@/frontend/dist

# TTS (set CCPIPE_TTS=off to disable)
Environment=CCPIPE_TTS=kokoro
Environment=CCPIPE_KOKORO_URL=http://localhost:8880
Environment=CCPIPE_TTS_VOICE=bf_emma

# Authentication is ALWAYS ON. ccpipe generates a random password on
# first start, argon2id-hashes it into ~/.local/state/ccpipe/credentials,
# and writes the plaintext (once) into a sidecar ~/.local/state/ccpipe/
# initial_password.txt (mode 0400). Read it, then `shred -u` the
# sidecar.
#
# To pin the credentials in advance (e.g. provisioning), use a drop-in
# with mode 0600 — NOT this version-controlled file:
#     mkdir -p ~/.config/systemd/user/ccpipe.service.d/
#     printf "[Service]\nEnvironment=CCPIPE_AUTH_PASSWORD=hunter2\n" \
#         > ~/.config/systemd/user/ccpipe.service.d/auth.conf
#     chmod 600 ~/.config/systemd/user/ccpipe.service.d/auth.conf
#     systemctl --user daemon-reload && systemctl --user restart ccpipe

Restart=on-failure
RestartSec=3

# StandardOutput/Error go to journal by default; view with:
#   journalctl --user -u ccpipe -f

[Install]
WantedBy=default.target
