# systemd/ccpipe.service.in
[Unit]
Description=ccpipe — Claude Code remote web interface
After=default.target ccpipe-virtual-mic.service
# Soft dep: try to bring up the virtual mic so /voice works. If it can't
# (e.g. Pulse not running), ccpipe still starts — text + TTS playback
# don't depend on it.
Wants=ccpipe-virtual-mic.service

[Service]
Type=exec
# @REPO_ROOT@ is substituted by scripts/install.sh at install time with
# the absolute path of the checkout being installed. The rendered unit
# lives in ~/.config/systemd/user/ccpipe.service; the source template
# stays in the repo.
WorkingDirectory=@REPO_ROOT@/backend
# Default ExecStart binds 0.0.0.0 so an off-host reverse proxy (nginx)
# can reach the backend over the LAN.
#
# `--proxy-headers --forwarded-allow-ips=127.0.0.1` is included by
# default for the **co-located** nginx case (proxy on the same host).
# Without these flags every request appears to come from the proxy's
# loopback IP and the per-IP login throttle collapses into a single
# global bucket. With them, uvicorn substitutes X-Forwarded-For when
# the immediate peer is in the allow-list.
#
# For an **off-host** nginx, override via a systemd drop-in that
# replaces ExecStart with the proxy's LAN IP:
#     mkdir -p ~/.config/systemd/user/ccpipe.service.d/
#     cat > ~/.config/systemd/user/ccpipe.service.d/proxy.conf <<'EOF'
#     [Service]
#     ExecStart=
#     ExecStart=@REPO_ROOT@/backend/.venv/bin/uvicorn ccpipe.main:app \
#         --host 0.0.0.0 --port 8080 \
#         --timeout-keep-alive 5 --limit-concurrency 200 \
#         --proxy-headers --forwarded-allow-ips=<nginx-LAN-ip>
#     EOF
# See README §"Reverse proxy".
#
# Slowloris hardening (pass-3 review #19): cap simultaneous in-flight
# requests so a single source can't saturate the worker by holding
# 500+ slow connections, and shorten the keep-alive idle window so an
# idle TCP doesn't pin a socket for the default 5 s × N connections.
# Operator-side nginx still controls client_body_timeout /
# client_header_timeout independently — this is the app's own bound.
ExecStart=@REPO_ROOT@/backend/.venv/bin/uvicorn ccpipe.main:app \
    --host 0.0.0.0 --port 8080 \
    --timeout-keep-alive 5 \
    --limit-concurrency 200 \
    --proxy-headers --forwarded-allow-ips=127.0.0.1

# Frontend bundle served by FastAPI when this exists
Environment=CCPIPE_FRONTEND_DIST=@REPO_ROOT@/frontend/dist

# TTS (set CCPIPE_TTS=off to disable)
Environment=CCPIPE_TTS=kokoro
Environment=CCPIPE_KOKORO_URL=http://localhost:8880
Environment=CCPIPE_TTS_VOICE=bf_emma

# Authentication is ALWAYS ON. ccpipe generates a random password on
# first start, argon2id-hashes it into ~/.local/state/ccpipe/credentials,
# and writes the plaintext (once) into a sidecar ~/.local/state/ccpipe/
# initial_password.txt (mode 0400). Read it, then `shred -u` the
# sidecar.
#
# To pin the *initial* credentials (e.g. provisioning), use a drop-in
# with mode 0600 — NOT this version-controlled file:
#     mkdir -p ~/.config/systemd/user/ccpipe.service.d/
#     printf "[Service]\nEnvironment=CCPIPE_AUTH_PASSWORD=hunter2\n" \
#         > ~/.config/systemd/user/ccpipe.service.d/auth.conf
#     chmod 600 ~/.config/systemd/user/ccpipe.service.d/auth.conf
#     systemctl --user daemon-reload && systemctl --user restart ccpipe
#
# CCPIPE_AUTH_PASSWORD is **bootstrap-only**: it seeds the credentials
# file when the file does not yet exist, then the file is authoritative.
# Later changes happen via Settings → Account in the UI. To re-seed
# from env after the file exists, delete the file
# (`rm ~/.local/state/ccpipe/credentials`) and restart ccpipe.

Restart=on-failure
RestartSec=3

# ─── Sandboxing — intentionally NOT enabled ──────────────────────────
# An earlier revision of this file applied the standard hardening
# stack (ProtectSystem=strict, SystemCallFilter=@system-service,
# RestrictNamespaces, etc.). For this workload it makes things WORSE,
# not better:
#   * ProtectSystem=strict mounts /tmp read-only, so tmux can't create
#     /tmp/tmux-$UID for its socket — every sticky-restore + new-
#     session 500s after a clean reboot (the dir only "worked" before
#     reboot because an interactive shell had pre-created it).
#   * SystemCallFilter=@system-service with @resources in the deny set
#     SIGSYS-kills the Claude Code CLI at startup: claude calls
#     sched_setscheduler(2) via pthread_setschedparam() to set thread
#     priorities, and that syscall lives in @resources. Removing
#     @resources from the deny set is not enough — claude's exact
#     syscall surface is opaque and changes between releases, so any
#     allow-list filter is a perpetual whack-a-mole that ends with
#     claude dying mid-session on a stranger's upgrade.
# The auth gate is the security boundary that matters; sandboxing the
# service process buys little when the whole point of the service is
# to spawn an arbitrary user-trusted Claude session anyway.

# StandardOutput/Error go to journal by default; view with:
#   journalctl --user -u ccpipe -f

[Install]
WantedBy=default.target
